# crontab
0 0 5 ? MON PATH=/usr/local/bin/ python /app/main.py > /proc/1/fd/1 2>/proc/1/fd/2
